# bind-stats.cron
#!/bin/bash
#
STATUSTMPFILE="/tmp/rndcstatus.tmp"
/usr/sbin/rndc status > ${STATUSTMPFILE}
