# bind-stats.cron
#!/bin/bash
#
STATSMPFILE="/tmp/rndcstatus.tmp"
/usr/sbin/rndc status > ${STATSMPFILE}
