# bind-stats.cron.d
*/5 * * * * root /usr/local/bin/bind-stats.cron
